# crontab
PATH="/usr/local/bin:/usr/bin:/bin"
0 2 * * * bash /server/scripts/iniciar.sh > /var/log/cron.log 2>/var/log/cron.log
# Cada día a las 2:00